# Distribution/pkg_scripts/postinstall
#!/bin/bash

launchctl unload /Library/LaunchDaemons/io.easylogin.EasyLoginDB.plist
launchctl load /Library/LaunchDaemons/io.easylogin.EasyLoginDB.plist

launchctl unload /Library/LaunchDaemons/io.easylogin.EasyLoginAgent.plist
launchctl load /Library/LaunchDaemons/io.easylogin.EasyLoginAgent.plist

exit 0
